# ur-last-fm/templates/ur-last-fm.timer.j2
[Unit]
Description=run ur-last-fm every hour

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
